# overlay/ipq40xx/etc/systemd/system/hostapd.service.d/01-ipq40xx.conf
[Unit]
# Wait for the kernel to initialize wireless hardware
BindsTo=sys-subsystem-net-devices-wl24g.device
After=sys-subsystem-net-devices-wl24g.device
