# overlay/tn48m/etc/systemd/system/fan.service
[Unit]
Description=FAN configuration service
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/fan

[Install]
WantedBy=multi-user.target