# overlay/cn9130-crb/etc/systemd/system/upgrade-prepare.service
[Unit]
Description=System update prepare
DefaultDependencies=no
ConditionPathExists=/usr/upgrade/switch-root

[Service]
Type=oneshot
ExecStart=/usr/upgrade/switch-root