# overlay/ipq40xx/etc/systemd/system/wifi-cal.service
[Unit]
Description=Extract Wi-Fi calibration data
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/sbin/wifi-cal

[Install]
WantedBy=sysinit.target
